# docker/php/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
directory = /etc/cron.d/
command = /usr/sbin/cron -f
autorestart = true

[program:php]
command = docker-php-entrypoint php-fpm